# passthrough-backups-igpu/backup-rules/10-vfio-pci.rules
SUBSYSTEM=="pci", ATTRS{idVendor}=="1002", ATTRS{idDevice}=="164e", RUN+="/usr/bin/vfio-bind"
